# scripts/clipity-web.service
[Unit]
Description=Clipity Web Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/clipity
Environment=NODE_ENV=production
EnvironmentFile=/opt/clipity/.env
ExecStart=/home/ubuntu/.bun/bin/bun run src/web/server.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clipity-web

[Install]
WantedBy=multi-user.target
